# sv/ocp/run
#!/bin/sh

[ ! -f /tmp/startup.lock ] && sleep 10 && exit 0;

exec /usr/local/bin/ocp 
